# cron.daily/directory-sweeper
#!/bin/sh 
[ -x "$HOME/bin/directory-sweeper" ] && "$HOME/bin/directory-sweeper"
